# crontab
* * * * * date >> /var/log/cron.log 2>&1
0 9 * * * /speedtest.sh hostname duration token user
0 0 * * * /alt_speedtest.sh iperf_hostname iperf_hostname
